# audio-rtsp.service
[Unit]
Description=Audio RTSP Streaming Service
After=network.target mediamtx.service
Requires=mediamtx.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/bin
ExecStart=/usr/local/bin/startmic.sh
Restart=always
RestartSec=10
# Give the service time to properly start all streams
TimeoutStartSec=30
# Set resource limits to ensure stability
LimitNOFILE=65536

# Make sure the process group is killed when the service is stopped
KillMode=mixed
KillSignal=SIGTERM

# Ensure environment is properly set up
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
