# config/mediamtx-audio.service
# LyreBirdAudio Stream Manager systemd Service
# Manages FFmpeg audio streams for MediaMTX RTSP server
#
# Installation:
#   sudo cp config/mediamtx-audio.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable mediamtx-audio
#   sudo systemctl start mediamtx-audio
#
# Note: Requires mediamtx.service to be running first

[Unit]
Description=LyreBirdAudio Stream Manager
Documentation=https://github.com/tomtom215/LyreBirdAudio
After=network.target mediamtx.service sound.target
Wants=mediamtx.service
PartOf=mediamtx.service

# Wait for audio subsystem
After=alsa-restore.service
Wants=alsa-restore.service

# Limit restart frequency (max 5 restarts in 300 seconds). These belong in
# [Unit]; systemd IGNORES StartLimitIntervalSec in [Service], so the rate limit
# there would silently not apply.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
# Type=forking (NOT simple): `start` daemonizes the FFmpeg wrappers and returns,
# writing the PID file below. With Type=simple systemd would treat that return as
# the service exiting and, with Restart=always, restart it in a loop.
Type=forking
PIDFile=/run/mediamtx-audio.pid
ExecStart=/usr/local/bin/lyrebird-stream-manager.sh start
ExecStop=/usr/local/bin/lyrebird-stream-manager.sh stop
ExecReload=/usr/local/bin/lyrebird-stream-manager.sh reload

# Working directory
WorkingDirectory=/var/lib/lyrebird

# User/Group (uncomment to run as non-root)
# User=mediamtx
# Group=mediamtx

# Restart policy for 24/7 operation
Restart=always
RestartSec=10
RestartPreventExitStatus=SIGTERM SIGINT

# NOTE: no WatchdogSec. The stream manager can feed systemd's watchdog only via
# the `systemd-notify` HELPER (a child process), which systemd rejects under the
# default NotifyAccess=main -- so a WatchdogSec would receive no keep-alive and
# kill the service every interval. LyreBird's own heartbeat, the optional
# hardware watchdog, and the cron health monitor provide recovery instead. To use
# the systemd watchdog anyway, add `NotifyAccess=all`, set WatchdogSec well above
# LYREBIRD_WATCHDOG_INTERVAL, and validate on your hardware first.

# Give time for graceful shutdown
TimeoutStopSec=30

# Environment. LYREBIRD_WATCHDOG_* drives LyreBird's OWN heartbeat / hardware
# watchdog (independent of systemd's WatchdogSec, which is intentionally unset
# above).
Environment=LYREBIRD_WATCHDOG_ENABLED=true
Environment=LYREBIRD_WATCHDOG_INTERVAL=30

# Optional: increase file descriptor limit for many streams
LimitNOFILE=65535

# Security hardening
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/log /var/lib/lyrebird /etc/mediamtx /run
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits for production stability
MemoryMax=1G
MemoryHigh=800M

# Logging
StandardOutput=append:/var/log/lyrebird-stream-manager.log
StandardError=append:/var/log/lyrebird-stream-manager.log
SyslogIdentifier=lyrebird-audio

[Install]
WantedBy=multi-user.target
